# edgefusion.service
# Replace __EDGEFUSION_USER__ and __EDGEFUSION_PROJECT_DIR__ before installing.
[Unit]
Description=EdgeFusion Backend
After=network.target

[Service]
Type=simple
User=__EDGEFUSION_USER__
Group=__EDGEFUSION_USER__
WorkingDirectory=__EDGEFUSION_PROJECT_DIR__
Environment=PYTHONUNBUFFERED=1
ExecStart=__EDGEFUSION_PROJECT_DIR__/run.sh
Restart=always
RestartSec=5
KillSignal=SIGINT
TimeoutStopSec=20

[Install]
WantedBy=multi-user.target
